# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true

[program:laravel-server]
command=php -S 0.0.0.0:3000 -t public
autostart=true
autorestart=true

